# dist/libcec-daemon.service
[Service]
Type=simple
ExecStart=/usr/bin/libcec-daemon

